# nra-whoisd.service
[Unit]
Description=Python Service


[Service]
ExecStart=/usr/local/bin/nra-whoisd.py
